# switchboard/rootfs/etc/s6-overlay/s6-rc.d/init-switchboard/up
/bin/sh -c "/usr/bin/switchboard-config && { ( /usr/bin/switchboard-mwi replay >/dev/null 2>&1 ) & }"
