# Switchboard/switchboard/rootfs/etc/s6-overlay/s6-rc.d/init-switchboard/up
#!/command/with-contenv bashio
# Render Asterisk configuration from the add-on options, then prepare runtime
# directories. Runs once, before Asterisk and the web UI start.
/usr/bin/switchboard-config
