# switchboard/rootfs/etc/s6-overlay/s6-rc.d/console-web/run
#!/command/with-contenv bashio
# Console web terminal — serves the operator console's TUI in the browser via
# xterm.js, bridging WebSocket <-> the telnet console on 127.0.0.1:2300. Add it
# to the Home Assistant sidebar with a panel_iframe pointing at this LAN port.
# host_network exposes the port directly on the LAN.
# Idle ONLY on an explicit `false`. `bashio::config.true` also returns false for an
# EMPTY/transient read — bashio can momentarily read blank options during a config
# reload / at boot — which would then `exec sleep infinity` and PERMANENTLY idle an
# ENABLED service (s6 sees an idle process as "successfully started" and never
# restarts it). Both console flags default enabled, so key off the literal value:
# empty/null (transient) ⇒ run.
if [ "$(bashio::config 'console_web_enabled')" = "false" ]; then
    bashio::log.info "Console web terminal disabled (console_web_enabled: false); idling."
    exec sleep infinity
fi

# The web terminal is a front-end for the operator console; if that's off there
# is nothing to bridge to, so idle instead of serving a dead page.
if [ "$(bashio::config 'console_enabled')" = "false" ]; then
    bashio::log.info "Operator console disabled (console_enabled: false); console web terminal idling."
    exec sleep infinity
fi

PORT="$(bashio::config 'console_web_port')"
# bashio can hand back an empty string / "null" during a config reload; keep the
# port sane so the log line and CONSOLE_WEB_PORT export are never empty (server.py
# also guards this, but a blank ":${PORT}" in the log is just misleading).
if [ -z "${PORT}" ] || [ "${PORT}" = "null" ]; then
    PORT="8100"
fi
bashio::log.info "Starting Switchboard console web terminal on :${PORT} ..."
bashio::log.notice "The console web terminal is UNAUTHENTICATED on the LAN and drives the operator console (ring/connect/hang up). It exposes the same surface as the telnet console on port 2300. Set console_web_enabled:false to turn it off."

export CONSOLE_WEB_PORT="${PORT}"
# Bind address: explicit console_web_bind wins; otherwise follow console_bind so
# locking the telnet console to loopback also locks the web terminal.
BIND="$(bashio::config 'console_web_bind')"
if [ -z "${BIND}" ] || [ "${BIND}" = "null" ]; then
    BIND="$(bashio::config 'console_bind')"
fi
if [ -z "${BIND}" ] || [ "${BIND}" = "null" ]; then
    BIND="0.0.0.0"
fi
export CONSOLE_WEB_BIND="${BIND}"
# Bridge to the operator console on loopback; mirror its configured port so a
# non-default console_port still works.
TARGET_PORT="$(bashio::config 'console_port')"
if [ -z "${TARGET_PORT}" ] || [ "${TARGET_PORT}" = "null" ]; then
    TARGET_PORT="2300"
fi
export CONSOLE_WEB_TARGET_PORT="${TARGET_PORT}"
# Timezone backstop, matching the other services (the page itself shows no
# clock, but keep the process TZ consistent).
export TZ="$(cat /run/switchboard/tz 2>/dev/null || echo UTC)"

cd /usr/share/switchboard/console-web || exit 1
exec python3 server.py
