# switchboard/rootfs/etc/s6-overlay/s6-rc.d/devhealth/run
#!/command/with-contenv bashio
# Proactive device-health monitor — polls the WP826 cordless's own HTTP API for the
# signals Asterisk can't see (battery %, Wi-Fi RSSI, per-call MOS) and derives GXW
# gateway health from rtpmon's rollup, alerting on a SINGLE critical device (the
# alarm cordless; the whole gateway) that the fleet-outage gate would miss.
# Idle only on an explicit `false` — an empty/transient bashio read at boot must
# not permanently idle this enabled-by-default monitor (s6 never restarts an idle).
if [ "$(bashio::config 'device_health_enabled')" = "false" ]; then
    bashio::log.info "Device-health monitor disabled (device_health_enabled: false); idling."
    exec sleep infinity
fi

bashio::log.info "Starting Switchboard device-health monitor ..."

export DEVICE_HEALTH_INTERVAL="$(bashio::config 'device_health_interval')"
export CORDLESS_IP="$(bashio::config 'cordless_ip')"
export CORDLESS_PASSWORD="$(bashio::config 'cordless_password')"
export GATEWAY_PORTS="$(bashio::config 'gateway_ports')"
export CORDLESS_BATTERY_CRIT="$(bashio::config 'cordless_battery_crit_pct')"
export CORDLESS_BATTERY_WARN="$(bashio::config 'cordless_battery_warn_pct')"
export CORDLESS_WIFI_MIN="$(bashio::config 'cordless_wifi_min_signal')"
if bashio::config.true 'device_health_alerts'; then
    export DEVICE_HEALTH_ALERTS=true
else
    export DEVICE_HEALTH_ALERTS=false
fi

exec python3 /usr/share/switchboard/devhealth/poller.py
